# snapserver/rootfs/etc/services.d/snapcast/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Runs the Snapcast Server
# ==============================================================================
declare -a options

bashio::log.info "Starting Snapcast Server..."

/usr/bin/snapserver -c /etc/snapcast/snapserver.conf